# dev-php/hhvm/files/hhvm.initd-2
#!/sbin/runscript
# Copyright 1999-2013 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

description="Virtual Machine, Runtime, and JIT for PHP"

HHVM_CONFIG="/etc/hhvm/${RC_SVCNAME}.hdf"

HHVM_LIB=/usr/lib/hhvm/lib
PIDDIR="/var/run/hhvm"
PIDFILE="${PIDDIR}/${RC_SVCNAME}.pid"
LOGDIR="/var/log/hhvm"
LOGFILE="${LOGDIR}/${RC_SVCNAME}-server.log"
HHVM_STOP_TIMEOUT=${HHVM_STOP_TIMEOUT:-15}
HHVM_UMASK=${HHVM_UMASK:-$(umask)}

depend() {
	need net localmount
	use dns logger
	after bootmisc
}

checkconfig() {
	if [ ! -f "${HHVM_CONFIG}" ] ; then
		eerror "Please create ${HHVM_CONFIG}"
		eerror "Sample conf: /etc/hhvm/config.hdf.dist"
		return 1
	fi
	return 0
}

start() {
	checkconfig || return $?

	checkpath -d --owner hhvm:hhvm --mode 0755 ${PIDDIR}
	checkpath -d --owner hhvm:hhvm --mode 0755 ${LOGDIR}

	ebegin "Starting HipHop VM (${RC_SVCNAME})"
	LD_LIBRARY_PATH="${HHVM_LIB}" start-stop-daemon --start --background \
		--exec "${HHVM_BINARY}" \
		--make-pidfile --pidfile "${PIDFILE}" \
		--stdout "${LOGFILE}" \
		--stderr "${LOGFILE}" \
		--env HPHP_HOME="${HHVM_LIB}/" \
        --umask "${HHVM_UMASK}" \
		-- -m daemon -u hhvm -c "${HHVM_CONFIG}"
	eend $? "Failed to start HipHop VM (${RC_SVCNAME})"
}

stop() {
	ebegin "Stopping HipHop VM (${RC_SVCNAME})"
	start-stop-daemon --stop --pidfile "${PIDFILE}" --retry "${HHVM_STOP_TIMEOUT}"
	eend $? "Failed to stop HipHop VM (${RC_SVCNAME})"
}
